# container/php/supervisor.conf
[program:caraudio]
process_name=%(program_name)s
command=php /var/www/html/caraudio/artisan queue:listen
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/caraudio.log
stopwaitsecs=3600
